# contrib/systemd/prmana-agent.service
[Unit]
Description=prmana Authentication Agent
Documentation=https://github.com/prodnull/unix-oidc
# Require the socket unit so systemd sets up LISTEN_FDS before starting the daemon.
Requires=prmana-agent.socket
After=prmana-agent.socket network-online.target

[Service]
# Type=notify: systemd waits for READY=1 sd-notify before marking the service active.
# The daemon sends READY=1 after socket acquisition, config validation, and an
# initial best-effort JWKS prefetch.
Type=notify
ExecStart=/usr/bin/prmana-agent serve
Restart=on-failure
RestartSec=5s

# --- Hardening directives ---
# Prevents the process from gaining new privileges via setuid/setgid or file capabilities.
NoNewPrivileges=yes

# Mount the entire file system read-only except /dev, /proc, /sys and
# any paths listed in ReadWritePaths below.
ProtectSystem=strict

# Allow reading from $HOME but block writes (private keys are in keyring or
# ReadWritePaths below; DPoP signing never writes to $HOME directly).
ProtectHome=read-only

# Required read-write paths under ProtectSystem=strict:
# - file storage fallback writes credentials here when no keyring is available
# - config file is read only, but directories must be traversable
ReadWritePaths=%h/.local/share/prmana %h/.config/prmana

# Prevent the daemon from executing writable memory (blocks JIT and most shellcode).
# If using a YubiKey/PKCS#11 signer that loads a vendor .so, uncomment:
# MemoryDenyWriteExecute=no
MemoryDenyWriteExecute=yes

# Block creation of new kernel namespaces (unshare, clone with namespace flags).
RestrictNamespaces=yes

# Prevent the daemon from setting SCHED_FIFO or SCHED_RR (real-time scheduling).
RestrictRealtime=yes

# Prevent setuid/setgid file execution inside this unit's filesystem view.
RestrictSUIDSGID=yes

# Provide a private /tmp and /var/tmp invisible to other processes.
PrivateTmp=yes

# Allow only syscalls required by a network service: standard I/O, memory management,
# signal handling, timers, and networking. This blocks exotic syscalls used by exploits.
# @system-service covers the base set; @network-io adds accept/bind/connect/send/recv.
SystemCallFilter=@system-service @network-io

[Install]
WantedBy=default.target
